# bm/battery-monitor.service
[Unit]
Description="Battery life"

[Service]
ExecStart=/usr/local/bin/battery-monitor

[Install]
WantedBy=default.target
