# bm/battery-monitor.timer
[Unit]
Description=Battery monitor timer

[Timer]
OnUnitActiveSec=30m
Unit=battery-monitor.service

[Install]
WantedBy=timers.target

